# update-hosts.timer
[Unit]
Description=Veckovis uppdatering av hosts-filen

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
