# update-hosts.timer
[Unit]
Description=Timer för uppdatering av hosts-filen

[Timer]
OnBootSec=15min
OnUnitActiveSec=6h

[Install]
WantedBy=timers.target
