# update-hosts.timer
[Unit]
Description=Update hosts file

[Timer]
OnBootSec=15min
OnUnitActiveSec=6h

[Install]
WantedBy=timers.target
